# launchd/com.terrancehan.quant-weekly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly REVIEW job: aggregates the week's daily-run JSONs, emails a
  summary. Observation-only — no parameter changes. Runs Friday 16:30
  LOCAL TIME (30 min after close).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.terrancehan.quant-weekly-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd '/Users/terrancehan/Claude trader/quant' &amp;&amp; /opt/homebrew/bin/uv run quant-weekly-review</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>16</integer>
        <key>Minute</key><integer>30</integer>
        <key>Weekday</key><integer>5</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/terrancehan/Claude trader/quant</string>

    <key>StandardOutPath</key>
    <string>/Users/terrancehan/Claude trader/quant/data/agent/launchd-logs/weekly-review.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/terrancehan/Claude trader/quant/data/agent/launchd-logs/weekly-review.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
