# launchd/com.terrancehan.quant-daily-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily REPORT job: loads today's persisted daily-run JSON, renders the
  markdown summary, emails it to REPORT_TO from .env. Runs at 16:05
  LOCAL TIME, Monday through Friday (5 minutes after NYSE close).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.terrancehan.quant-daily-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd '/Users/terrancehan/Claude trader/quant' &amp;&amp; /opt/homebrew/bin/uv run quant-daily-report</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>1</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>2</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>4</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>5</integer><key>Weekday</key><integer>5</integer></dict>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/terrancehan/Claude trader/quant</string>

    <key>StandardOutPath</key>
    <string>/Users/terrancehan/Claude trader/quant/data/agent/launchd-logs/daily-report.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/terrancehan/Claude trader/quant/data/agent/launchd-logs/daily-report.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
