# launchd/com.terrancehan.quant-daily-trade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily TRADE job: cancels open orders, computes today's target weights
  from yesterday's close, submits market entries with atomic 5% stop-losses.
  Runs at 21:35 CHINA STANDARD TIME, Monday through Friday.

  Mac is on CST (+0800). 21:35 CST same day maps to:
    - 09:35 EDT (summer, US-DST): 5 min after market open. Ideal.
    - 08:35 EST (winter, no-DST): 55 min before market open. Orders
      submitted then queue at Alpaca and fill at the opening cross —
      acceptable (we only submit market orders + OTO stops, all of
      which fill at open or near-open prices either way).

  TZ note for the next maintainer: this plist's schedule is hardcoded
  for CST. If you move the Mac to a different timezone, edit the
  Hour/Weekday values. If you want to eliminate the EDT/EST 1-hour
  drift entirely, set the Mac TZ to America/New_York and use 09:35
  / Mon-Fri instead.

  Logs go to data/agent/launchd-logs/.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.terrancehan.quant-daily-trade</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd '/Users/terrancehan/Claude trader/quant' &amp;&amp; /opt/homebrew/bin/uv run quant-daily-trade</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer><key>Weekday</key><integer>1</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer><key>Weekday</key><integer>2</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer><key>Weekday</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer><key>Weekday</key><integer>4</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer><key>Weekday</key><integer>5</integer></dict>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/terrancehan/Claude trader/quant</string>

    <key>StandardOutPath</key>
    <string>/Users/terrancehan/Claude trader/quant/data/agent/launchd-logs/daily-trade.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/terrancehan/Claude trader/quant/data/agent/launchd-logs/daily-trade.err</string>

    <!-- RunAtLoad false: don't fire immediately on install, only at the
         scheduled time. The operator can manually invoke `uv run
         quant-daily-trade --dry-run` first to sanity-check. -->
    <key>RunAtLoad</key>
    <false/>

    <!-- KeepAlive auto-retry: if the trade routine exits non-zero
         (e.g. DNS flake, VPN tunnel re-establishment, transient broker
         outage), launchd re-fires the entire job. Combined with
         ThrottleInterval=300, retries are spaced 5 minutes apart, so
         the loop tops out at ~12 retries per hour.

         The retry loop is bounded automatically:
           - On success: run_daily_trade saves today's run JSON; next
             fire detects it via the idempotency guard and exits clean.
           - On persistent failure: _past_trade_deadline() (15:35 ET)
             returns early with exit 0 once the trade window closes,
             ending the retry loop until tomorrow's StartCalendarInterval.

         Without KeepAlive, a one-time DNS hang at 09:35 ET = a lost
         trading day (June 2026 incident). With it, the same hang costs
         5 minutes at most. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
